# worker/patches/etc/profile.d/oraclejdk.sh
export J2SDKDIR=/opt/java/jdk1.8.0_33
export J2REDIR=/opt/java/jdk1.8.0_33/jre
export PATH=$PATH:/opt/jdk1.8.0_33/bin:/opt/java/jdk1.8.0_33/db/bin:/opt/java/jdk1.8.0_33/jre/bin
export JAVA_HOME=/opt/java/jdk1.8.0_33
export DERBY_HOME=/opt/java/jdk1.8.0_33/db
